# tests/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; the path to the socket file

[supervisord]
logfile=/tmp/supervisord.log ; main log file; default $CWD/supervisord.log
logfile_maxbytes=50MB        ; max main logfile bytes b4 rotation; default 50MB
logfile_backups=10           ; # of main logfile backups; 0 means none, default 10
loglevel=info                ; log level; default info; others: debug,warn,trace
pidfile=/tmp/supervisord.pid ; supervisord pidfile; default supervisord.pid
nodaemon=false               ; start in foreground if true; default false
silent=false                 ; no logs to stdout if true; default false
minfds=1024                  ; min. avail startup file descriptors; default 1024
minprocs=200                 ; min. avail process descriptors;default 200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket

[program:crond]
command = /usr/sbin/crond -n -m off
startsecs = 0
priority = 100
autostart = true
autorestart = true
stdout_logfile = /var/log/supervisor/crond.log
redirect_stderr = true

#[program:crond]
#priority = 100
#command = bash -c "while true; do sleep 0.1; [[ -e /var/run/crond.pid ]] || break; done && exec /usr/sbin/crond -m off -n"
#startsecs = 0
#autorestart = true
#redirect_stderr = true
#stdout_logfile = /var/log/cron
#stdout_events_enabled = true

[program:mysqld]
command = /usr/bin/pidproxy /var/run/mariadb/mariadb.pid /usr/bin/mysqld_safe
stdout_logfile = /var/log/supervisor/mysqld.log
redirect_stderr = true
user = mysql
autostart = true

[program:munged]
command= /usr/sbin/munged -F
user = munge
autostart = false

[program:mailserver]
command = /usr/local/sbin/mail-server.py -i 127.0.0.1 -p 25
priority = 3
autostart = true
stdout_logfile = /var/log/supervisor/mailserver.log
redirect_stderr = true

[program:slurmdbd]
command = /usr/sbin/slurmdbd -D
user = root
autostart = false

[program:slurmctld]
command = /usr/sbin/slurmctld -D
user = root
autostart = false

[program:slurmd]
command = /usr/sbin/slurmd -D
user = root
autostart = false
