# etc/cron.d/slurm-mail
*    *    *    *    *    root    /usr/bin/slurm-send-mail
